# bin/supervisor/queue-consumer/g3-queue-worker.conf
[group:g3-queue-workers]
programs=g3-queue-default,g3-queue-high

[program:g3-queue-default]
; 默认队列 - 处理一般优先级任务（包括邮件等所有任务类型）
command=php %(here)s/../queue-worker.php --daemon --verbose --queue=default --sleep=3 --tries=3
directory=/var/www/html
user=www-data
group=www-data

; 默认队列：2个并行进程
numprocs=2
process_name=%(program_name)s_%(process_num)02d

priority=100
autostart=true
autorestart=true
startsecs=10
startretries=3
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=10

[program:g3-queue-high]
; 高优先级队列 - 处理紧急任务
command=php %(here)s/../queue-worker.php --daemon --verbose --queue=high --sleep=1 --tries=5 --timeout=60
directory=/var/www/html
user=www-data
group=www-data

; 高优先级队列：1个进程 
numprocs=1
process_name=%(program_name)s_%(process_num)02d

priority=200
autostart=true
autorestart=true
startsecs=5
startretries=5
exitcodes=0,2
stopsignal=TERM
stopwaitsecs=15

; 通用日志配置 - 所有队列共享
stdout_logfile=/var/log/supervisor/g3-queue-%(program_name)s.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
stderr_logfile=/var/log/supervisor/g3-queue-%(program_name)s-error.log
stderr_logfile_maxbytes=10MB
stderr_logfile_backups=5
redirect_stderr=false

; 通用环境变量
environment=
    PATH="/usr/local/bin:/usr/bin:/bin",
    HOME="/var/www",
    USER="www-data",
    LANG="en_US.UTF-8",
    LC_ALL="en_US.UTF-8"

; 安全配置
umask=022